# systemd/gsheets-radius-sync.timer
[Unit]
Description=Daily synchronize FreeRADIUS database with Google Sheets

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
