# systemd/gsheets-dhcp-gen.service
[Unit]
Description=Update DHCPd config from GSheets

[Service]
Type=oneshot
ExecStart=/usr/local/bin/gsheets-dhcp-gen.sh

