# systemd/gsheets-vlan-gen.timer
[Unit]
Description=Daily update VLAN config from GSheets

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
