# systemd/gsheets-vlan-gen.service
[Unit]
Description=Update VLAN config from GSheets

[Service]
Type=oneshot
ExecStart=/usr/local/bin/gsheets-vlan-gen.sh

